# unit-section/fail-5s.service
[Unit]
Description=fail in 5s

[Service]
Type=simple
ExecStart=/usr/local/bin/wait.sh 2-fail-5s 5
