# unit-section/fail-5s.service
[Unit]
Description=fail in 5s
OnFailure=status_email_user@%n.service

[Service]
Type=simple
ExecStart=/usr/local/bin/wait.sh 2-fail-5s 5
