# unit-section/fail-immediate.service
[Unit]
Description=fail immediate

[Service]
Type=simple
ExecStart=/usr/local/bin/wait.sh 3-fail-immediate 0
